# Sources/Initrd-tcos1.0/etc/udev/rules.d/80-compcache.rules
KERNEL=="ramzswap0", ACTION=="add", \
  RUN+="/sbin/swapon -p 100 /dev/ramzswap0 2>/dev/null"
